# launchd/com.llmeter.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.llmeter.monitor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/path/to/llmeter/.venv/bin/python</string>
    <string>-m</string>
    <string>llmeter</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/path/to/llmeter</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/john/.llmeter/logs/llmeter.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/john/.llmeter/logs/llmeter.log</string>
</dict>
</plist>
